# mpc/root/etc/rc.local
#!/bin/bash
mkdir -p /config/.java
chown -R nobody:users /config
rm -f /tmp/.X1-lock
exit 0

